# decker-engine.service
[Unit]
Description=Decker Engine

[Service]
Type=forking
Environment="DECKER_BASE_URL=https://tramberend.beuth-hochschule.de/decker" "DECKER_ROOT_DIR=/var/local/decker"
ExecStart=/usr/bin/daemonize \
  -c /var/local/decker
  -u decker
  /usr/local/bin/decker-engine

[Install]
WantedBy=multi-user.target
