# decker-engine.service
[Unit]
Description=Decker Engine

[Service]
Type=simple
WorkingDirectory=/var/local/decker
ExecStart=/usr/local/bin/decker-engine-exe
User=daemon
Group=daemon
Restart=always
RestartSec=1s
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=decker-engine

[Install]
WantedBy=multi-user.target

